# subprojects/metric-collection-daemon/phosphor-metric-collector.service.in
[Unit]
Description=BMC Metric Collection

[Service]
ExecStart=@bindir@/metric-collector
Type=dbus
BusName=xyz.openbmc_project.MetricCollector
SyslogIdentifier=phosphor-metric-collector

[Install]
WantedBy=multi-user.target
