# subprojects/bare-metal-host-monitor/host-gpio-monitor@.service.in
[Unit]
Description=gBMC bare-metal GPIO  monitoring for host %i
ConditionPathExists=/var/google/config-package/enable-bm.flag
After=phosphor-ipmi-host.service
After=xyz.openbmc_project.Chassis.Control.Power@%i.service
After=xyz.openbmc_project.State.Host@%i.service

[Service]
Type=notify
Restart=always
ExecStart=@@BIN@ host_gpio_monitor --host_label %i

[Install]
WantedBy=multi-user.target
