# acpi-power-state-daemon/acpi-power-state-daemon/acpi-power-state.service.in
[Unit]
Description=Receive ACPI power state over IPMI

[Service]
Type=simple
Restart=always
ExecStart=@@BIN@ acpi_power_stated

[Install]
WantedBy=multi-user.target
