# acpi-power-state-daemon/acpi-power-state-daemon/host-s5-state.target
[Unit]
Description=Host S5 state target
Conflicts=host-s0-state.target
