# subprojects/acpi-power-state-daemon/host-s0-state.target
[Unit]
Description=Host S0 state target
Conflicts=host-s5-state.target
