# subprojects/kcsbridge/src/kcsbridge@.service.in
[Unit]
Description=Google IPMI KCS Bridge

[Service]
Restart=on-failure
Type=notify
ExecStart=@BIN@ -c '%I'

[Install]
WantedBy=multi-user.target
