# subprojects/dhcp-done/dhcp-done.service.in
[Unit]
Description=gBMC DHCP Status Daemon

[Service]
Restart=on-failure
ExecStart=@@BIN@ dhcp-done

[Install]
WantedBy=multi-user.target
