# subprojects/dhcp-done/dhcp-done@.service.in
[Unit]
Description=gBMC DHCP Complete

[Service]
Restart=on-failure
ExecStart=@@BIN@ dhcp-done %I
