# subprojects/dhcp-done/dhcp-done.service.in
[Unit]
Description=gBMC DHCP Status Daemon

[Service]
Restart=on-failure
ExecStartPre=update-dhcp-status 'ONGOING' 'Waiting on dhcp process'
ExecStart=@@BIN@ dhcp-done

[Install]
WantedBy=multi-user.target
