# ncsid/src/nic-hostful@.target
[Unit]
Description=Target when nic enters mode with host
Conflicts=nic-hostless@%i.target
