# subprojects/btmanager/service_files/btmanager.service.in
[Unit]
Description=Host boot time manager
After=xyz.openbmc_project.Chassis.Control.Power.service

[Service]
Type=oneshot
Restart=on-failure
ExecStart=@BIN@

[Install]
WantedBy=multi-user.target
