# subprojects/btmanager/service_files/btmanager.service.in
[Unit]
Description=Host boot time manager

[Service]
Restart=on-failure
ExecStart=@BIN@
