# subprojects/ncsid/src/dhcp6@.service.in
[Unit]
Description=DHCPv6 Updater
Wants=mapper-wait@-xyz-openbmc_project-network-%i.service
After=mapper-wait@-xyz-openbmc_project-network-%i.service
Requisite=nic-hostless@%i.target
After=nic-hostless@%i.target
BindsTo=nic-hostless@%i.target
ConditionPathExists=!/run/dhcp6.done

[Service]
KillMode=mixed
Restart=on-failure
ExecStart=@/usr/bin/env udhcpc6 -f -q -O bootfile_url -O bootfile_param -i %I -s @SCRIPT@
SyslogIdentifier=dhcp6@%I
