# subprojects/ncsid/src/ncsid@.service.in
[Unit]
Description=Run Ncsid daemon
Wants=update-ra-gw@%i.service

[Service]
Restart=always
ExecStart=@@BIN@ ncsid %I
SyslogIdentifier=ncsid@%I

[Install]
WantedBy=multi-user.target
