# subprojects/ncsid/src/ncsid@.service.in
[Unit]
Description=Run Ncsid daemon
Wants=systemd-networkd.service
After=systemd-networkd.service

[Service]
Restart=always
ExecStart=@@BIN@ ncsid %I
SyslogIdentifier=ncsid@%I

[Install]
WantedBy=multi-user.target
