# subprojects/ncsid/src/update-ra-gw@.service.in
[Unit]
Description=RA Gateway Updater
Wants=mapper-wait@-xyz-openbmc_project-network-%i.service
After=mapper-wait@-xyz-openbmc_project-network-%i.service
After=ncsid@%i.service
StartLimitIntervalSec=1min
StartLimitBurst=5

[Service]
KillMode=mixed
Restart=on-failure
ExecStart=@@BIN@ update-ra-gw %I
SyslogIdentifier=update-ra-gw@%I
SuccessExitStatus=10
