# subprojects/ncsid/src/update-ra-neighbor@.timer
[Unit]
Description=Retry neighbor update periodically

[Timer]
OnUnitInactiveSec=1min
