# subprojects/ncsid/src/update-ra-neighbor@.timer
[Unit]
Description=Retry neighbor update periodically
BindsTo=ncsid@%i.service
After=ncsid@%i.service

[Timer]
OnUnitInactiveSec=1min
