# subprojects/ncsid/src/update-ra-neighbor@.service.in
[Unit]
Description=RA Neighbor Updater
Wants=mapper-wait@-xyz-openbmc_project-network-%i.service
After=mapper-wait@-xyz-openbmc_project-network-%i.service
StartLimitIntervalSec=1min
StartLimitBurst=5

[Service]
KillMode=mixed
Restart=on-failure
ExecStart=@@BIN@ update-ra-neighbor %I
SyslogIdentifier=update-ra-neighbor@%I
SuccessExitStatus=10
